# launchd/csu.autorelogin.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>cn.csu.autorelogin</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/your-name/Documents/CSU-student-WiFi/auto_relogin.py</string>
        <string>--config</string>
        <string>/Users/your-name/Library/Application Support/CSUStudentWiFi/config.toml</string>
        <string>--once</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/your-name/Documents/CSU-student-WiFi</string>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>18000</integer>
    <key>StandardOutPath</key>
    <string>/Users/your-name/Library/Application Support/CSUStudentWiFi/launchd.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/your-name/Library/Application Support/CSUStudentWiFi/launchd.stderr.log</string>
</dict>
</plist>
