# launchd/csu.autorelogin.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>cn.csu.autorelogin</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/zhangyu/Documents/my-project/CSU-student-WiFi/.venv/bin/python</string>
        <string>/Users/zhangyu/Documents/my-project/CSU-student-WiFi/auto_relogin.py</string>
        <string>--config</string>
        <string>/Users/zhangyu/Documents/my-project/CSU-student-WiFi/config.toml</string>
        <string>--once</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/zhangyu/Documents/my-project/CSU-student-WiFi</string>
    <key>RunAtLoad</key>
    <true/>
    <key>StartInterval</key>
    <integer>18000</integer>
    <key>StandardOutPath</key>
    <string>/tmp/csu_autorelogin.stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/tmp/csu_autorelogin.stderr.log</string>
</dict>
</plist>
